# src/etc/s6-overlay/s6-rc.d/init-start/up
# This file doesn't do anything, it's just starting point.
